# openclaw/launch_agents/ai.openclaw.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.watchdog</string>

    <key>Comment</key>
    <string>OpenClaw Watchdog (SSOT @ xmad-control/openclaw)</string>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>ThrottleInterval</key><integer>30</integer>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/ahmadabdullah/xmad-control/openclaw/scripts/watchdog.sh</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/ahmadabdullah/xmad-control/openclaw/logs/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ahmadabdullah/xmad-control/openclaw/logs/watchdog.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>NODE_ENV</key>
      <string>production</string>
      <key>HOME</key>
      <string>/Users/ahmadabdullah</string>
      <key>OPENCLAW_STATE_DIR</key>
      <string>/Users/ahmadabdullah/xmad-control/openclaw</string>
      <key>PATH</key>
      <string>/Users/ahmadabdullah/.bun/bin:/Users/ahmadabdullah/.nvm/versions/node/v22.16.0/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
